# bin/usbip/client/usbip-client.service
[Unit]
Description=usbip-client
After=network.target
Before=iobroker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c "usbip attach -r 192.168.1.41 -b $(usbip list -r 192.168.1.41 | grep '1-1.2' | cut -d: -f1 | sed -e 's/^[[:space:]]*//'); usbip attach -r 192.168.1.41 -b $(usbip list -r 192.168.1.41 | grep '1-1.3' | cut -d: -f1 | sed -e 's/^[[:space:]]*//');iob restart smartmeter.0; iob restart smartmeter.1"
ExecStop=/bin/sh -c "/usr/sbin/usbip detach --port=$(/usr/sbin/usbip port | grep '<Port in Use>' | sed -E 's/^Port ([0-9][0-9]).*/\1/');/usr/sbin/usbip detach --port=$(/usr/sbin/usbip port | grep '<Port in Use>' | sed -E 's/^Port ([0-9][0-9]).*/\1/'); iob stop smartmeter.0; iob stop smartmeter.1"

[Install]
WantedBy=multi-user.target 
