# bin/usbip/server/usbipd-bind@.service
# https://github.com/darkiop/doku-iobroker-ha
#
# /lib/systemd/system/usbipd-bind@.service
#

[Unit]
Description=usbipd binding %I
After=network-online.target usbipd.service
Wants=network-online.target
Requires=usbipd.service

[Service]
Type=simple
ExecStart=/usr/sbin/usbip --log bind -b %i
RemainAfterExit=yes
ExecStop=/usr/sbin/usbip --log unbind -b %i
Restart=on-failure

[Install]
WantedBy=multi-user.target