# bin/iobroker/socat/systemd/socat-loki-usb1.service
# /etc/systemd/system/socat-loki-usb1.service

[Unit]
Description=SOCAT Listen Service
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=iobroker
ExecStart=/bin/bash -c '/usr/bin/socat pty,link=/tmp/ttyCOM1 tcp:192.168.1.41:9991'
Restart=on-failure

[Install]
WantedBy=multi-user.target