# bin/iobroker/usbip/server/systemd/usbipd-bind@.service
# https://github.com/darkiop/doku-iobroker-ha
#
# /etc/systemd/system/usbipd-bind@.service
#
# systemctl enable usbipd-bind@1-1.2

[Unit]
Description=usbipd binding %I
After=network-online.target usbipd.service
Wants=network-online.target
Requires=usbipd.service

[Service]
Type=simple
ExecStart=/usr/sbin/usbip --log bind -b %i
RemainAfterExit=yes
ExecStop=/usr/sbin/usbip --log unbind -b %i
Restart=on-failure

[Install]
WantedBy=multi-user.target