# bin/proxmox/systemd/gstatus-status-to-json.timer
[Unit]
Description=gstatus to json timer
[Timer]
OnBootSec=1min
OnUnitInactiveSec=1min
Persistent=true
Unit=gstatus-status-to-json.service
[Install]
WantedBy=timers.target